# sys/cron/crontab
* * 01 * * /usr/bin/docker exec -it mail_opendkim_1 opendkim-rekey
